# dotfiles/theme-switcher/.config/systemd/user/theme-switcher.service
[Unit]
Description=Theme Switcher
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/theme-switcher.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
